# root/etc/service/nginx/run
#!/bin/sh
exec /usr/sbin/nginx -g "daemon off;" >>/log/nginx.log 2>&1